# scripts/metallb-monitor.service
[Unit]
Description=MetalLB LoadBalancer Monitor
After=k3s.service
Requires=k3s.service

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/monitor-metallb.sh --daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
